# private_dot_config/systemd/user/pi-maintenance-agent.service.tmpl
[Unit]
Description=Pi Maintenance Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
UMask=0077
EnvironmentFile=-%h/.config/dotfiles/pi-maintenance-agent.env
WorkingDirectory=%h/.local/share/pi-maintenance-agent
ExecStart=%h/.local/share/pi-maintenance-agent/bin/run-maintenance.sh
Nice=10
NoNewPrivileges=true
# Several filesystem/mount-namespace hardening knobs make OpenSSH see the
# packaged /etc/ssh/ssh_config.d systemd include as nobody-owned on Arch,
# which breaks non-interactive git/jj pushes from this service.
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native

[Install]
WantedBy=default.target
