# private_dot_config/systemd/user/pi-maintenance-agent.timer
[Unit]
Description=Run Pi Maintenance Agent on weekday evenings

[Timer]
OnCalendar=Mon..Fri 17:01 America/New_York
Persistent=true
Unit=pi-maintenance-agent.service

[Install]
WantedBy=timers.target
